# srv.service
[Unit]
Description=Kohlschwarz Think-Tank Apps Showcase

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/shelley-apps
ExecStart=/home/exedev/shelley-apps/server
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
